# deploy/systemd/jetsurf-api.service
[Unit]
Description=JetSurf leads API (Express + Telegram)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/jetsurf-wts/server
Environment=NODE_ENV=production
EnvironmentFile=-/var/www/jetsurf-wts/server/.env
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
